# sample_crontab/test.crontab

* * * * * echo "test `date`" >> /var/log/run.log
# ======= running in the same instance
# Run show date every 1 minute
* * * * * echo "output @ `date`" >> /var/log/run-cron_serv.log
# Save something every 2 minutes in an hour to target log file
*/2 * * * * echo "log @ `date`" >> /var/log/run2.log
